# deploy/ironclaw.service
[Unit]
Description=IronClaw AI Assistant
After=cloud-sql-proxy.service docker.service
Requires=cloud-sql-proxy.service

[Service]
Type=simple
EnvironmentFile=-/opt/ironclaw/service.env
ExecStartPre=-/usr/bin/docker rm -f ironclaw
ExecStart=/usr/bin/bash -lc '/usr/bin/docker run --rm \
  --name ironclaw \
  --env-file /opt/ironclaw/.env \
  --publish 127.0.0.1:${TITANCLAW_HOST_PORT:-3000}:3000 \
  ${TITANCLAW_IMAGE:?set TITANCLAW_IMAGE in /opt/ironclaw/service.env} \
  --no-onboard'
ExecStop=/usr/bin/docker stop ironclaw
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
